# channels/brain_loop/albunyaan-logo-sampler.timer
[Unit]
Description=Trigger Albunyaan logo sampler every 3 minutes
Documentation=file:///home/msa/Development/scripts/albunyaan/channels/brain_loop/logo_sampler.sh

[Timer]
# 3-minute cadence. Rule 10 (lessons.db) requires 5 rounds × 3 min =
# 12 min of observation before any mismatch flag. Tighten/loosen here
# if rule 10 evolves; no brain or watcher code has to change.
OnBootSec=90s
OnUnitActiveSec=3min

# Run the job at least every 4 min even if a prior run is still active.
AccuracySec=15s

Unit=albunyaan-logo-sampler.service
Persistent=true

[Install]
WantedBy=timers.target
