# channels/albunyaan-watcher.service
[Unit]
Description=Albunyaan reflex watcher (HLS freshness + system resources)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=msa
Group=msa
WorkingDirectory=/home/msa/Development/scripts/albunyaan/channels
RuntimeDirectory=albunyaan
# 0750: target_url files in cmd/ carry provider auth tokens (Xtream,
# seenshow, signed CDN URLs). Restrict world read while preserving group
# access for any future same-group sidecar tools.
RuntimeDirectoryMode=0750
ExecStartPre=/home/msa/Development/scripts/albunyaan/channels/reflex/preflight.sh
ExecStart=/home/msa/Development/scripts/albunyaan/channels/reflex_watcher.sh
Restart=always
RestartSec=5
Nice=10
CPUQuota=10%
MemoryMax=128M

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
